# .config/systemd/user/startapp.service
[Unit]
Description=Autostart launch hyperland

[Service]
ExecStart=/home/Archer/.dotfiles/script/start

[Install]
WantedBy=default.target

